# manageFiles/manageFiles.conf
[program:manageFiles]

command     = gunicorn -w4 -b127.0.0.1:8000 manageFiles.wsgi
directory   = /srv/manageFiles/manageFiles
startsecs=3
stopwaitsecs=0
autostart=true
autorestart=true


